# dpkg/prototype/DEBIAN/postinst
#!/bin/sh
update-rc.d hostname-setup start 01 S

