# dpkg/prototype/DEBIAN/postinst
#!/bin/sh
update-rc.d hostname-prompt start 01 S
update-rc.d hostname-persist start 01 S

