# dpkg/prototype/DEBIAN/prerm
#!/bin/sh
update-rc.d -f hostname-setup remove

