# dpkg/prototype/DEBIAN/prerm
#!/bin/sh
update-rc.d -f hostname-prompt remove
update-rc.d -f hostname-persist remove

